# voice-assistant.service
[Unit]
Description=Voice Assistant with Wake Word Detection
After=network.target sound.target pulseaudio.service
Wants=network-online.target

[Service]
Type=simple
User=knissel
Group=knissel
WorkingDirectory=/home/knissel/voice_assist
Environment=PATH=/home/knissel/voice_assist/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/knissel

# Main executable
ExecStart=/home/knissel/voice_assist/.venv/bin/python wakeword.py

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-assistant

# Resource limits for Pi 5
CPUQuota=80%
MemoryMax=768M
MemoryHigh=512M

# Ensure runtime dirs exist
ExecStartPre=/bin/mkdir -p /home/knissel/voice_assist/logs

# Security hardening (optional, may need adjustment)
# NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=read-only
# ReadWritePaths=/tmp /home/knissel/voice_assist/logs
# PrivateTmp=true

[Install]
WantedBy=multi-user.target
