# voice-assistant-ui.service
[Unit]
Description=Voice Assistant UI Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=knissel
Group=knissel
WorkingDirectory=/home/knissel/voice_assist
Environment=PATH=/home/knissel/voice_assist/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/knissel

ExecStart=/home/knissel/voice_assist/.venv/bin/python ui_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
